# ops/supervisor/wedigbio-ingest.conf
[program:wedigbio-ingest]
directory=/data/web/wedigbio-reports/current
command=/usr/bin/php artisan queue:work beanstalkd --queue=wedigbio-ingest --sleep=3 --tries=3 --timeout=120 --max-time=3600
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=60
stopasgroup=true
killasgroup=true
user=ubuntu
redirect_stderr=true
stdout_logfile=/data/web/wedigbio-reports/shared/storage/logs/supervisor-queue-worker.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
environment=APP_ENV="production"

